# supervisord.conf
[supervisord]
nodaemon=true

[program:fastapi]
command=bash ./supervisor/commands/fastapi.sh
autorestart=true
stderr_logfile=/var/log/fastapi.log