# supervisord.conf
[supervisord]
nodaemon=false

[program:fastapi]
command=zsh ./supervisor/commands/fastapi.sh
autorestart=true
stderr_logfile=./supervisor/logs/fastapi.log